# infra/module-manager-v2.service
[Unit]
Description=Module Manager V2 API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/module-manager-v2/v2-api
EnvironmentFile=/opt/module-manager-v2/.env
ExecStart=/opt/module-manager-v2/.venv/bin/python -m uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
User=modulemgr
Group=modulemgr

[Install]
WantedBy=multi-user.target
